# config/crontab
@reboot /usr/sbin/sshd -f /etc/ssh/sshd_config_2

